# dev-o-matic/interactive-image/docker/service/rabbitmq/run
#!/bin/sh

exec 2>&1
exec postgres -D /usr/local/pgsql/data >/var/log/postgres 2>&1